# com.intelligence-stack.provision-webui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.intelligence-stack.provision-webui</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yusuf/Projects/intelligence-stack/provision-webui.sh</string>
    </array>

    <!-- Run once at login. No KeepAlive — this is a one-shot provisioner. -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/provision-webui.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/provision-webui.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/yusuf</string>
    </dict>
</dict>
</plist>
